# etc/service/peerdocument/run
#!/bin/bash -e

cd /peerdocument

export DEBUG=peerdocument:*

exec chpst -u peerdocument:peerdocument /usr/local/bin/npm start 2>&1
